# DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable sdrivemontpn.timer
systemctl start sdrivemontpn.timer
